# gearman/worker_config/worker.conf
[program:mygearman]
command=php worker2.php
directory=/var/www/html/gearman/
autostart=true
autorestart=true
user=root
startsecs=0
